# root/etc/s6-overlay/s6-rc.d/svc-pas/run
#!/usr/bin/with-contenv bash

cd ${PAS_PATH} || exit

exec \
    s6-setuidgid abc ${PAS_PATH}/venv/bin/python3 ${PAS_PATH}/main.py;